# config/beamline_status_logger.target
[Unit]
Description=Log beamline attributes to a database
Requires=multi-user.target

[Install]
WantedBy=multi-user.target
