# config/beamline_status_logger.target
[Unit]
Description=Log beamline attributes to a database
After=network.target

[Install]
WantedBy=multi-user.target
