# config/beamline_status_logger@.service.in
[Unit]
Description=Log beamline attributes to a database
After=network.target
StartLimitIntervalSec=100
StartLimitBurst=5

[Service]
Type=simple
User=beamline-status-logger
ExecStart={script} {confdir}%i
Restart=on-failure
RestartSec=10

[Install]
WantedBy=beamline_status_logger.target
